# deployment/upstart.conf
start on started network
stop on stopping network

respawn

env AWS_ACCESS_KEY=""
env AWS_SECRET_KEY=""
env PORT=80
env NODE_ENV=production

exec node /opt/proxxy-node/index.js
